# automation/ansible/templates/job_scraper.timer.j2
[Unit]
Description=Daily job scraper timer

[Timer]
OnCalendar={{ job_scraper_run_time | default('*-*-* 08:00:00') }}
Persistent=true

[Install]
WantedBy=timers.target
