# automation/ansible/templates/minecraft-playerstate-backup.timer.j2
[Unit]
Description=Minecraft player-state backup timer

[Timer]
OnCalendar={{ minecraft_playerstate_backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target