# automation/ansible/templates/minecraft-dialogue-schedule.timer.j2
[Unit]
Description=Oakcross dialogue schedule timer

[Timer]
OnCalendar={{ minecraft_dialogue_schedule_calendar }}
Persistent=true

[Install]
WantedBy=timers.target